# usr/lib/systemd/system/custom_backup.service
[Unit]
Description=Full Backup to CHESS

[Service]
Type=simple
ExecStart=/usr/local/bin/system_backup.sh
